# etc/supervisord.conf
[unix_http_server]
file=%(here)s/../var/supervisor.sock
chmod=0600

[inet_http_server]
port = 127.0.0.1:4999

[supervisorctl]
serverurl=unix://%(here)s/../var/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
logfile=%(here)s/../var/supervisord.log
pidfile=%(here)s/../var/supervisord.pid
childlogdir=%(here)s/../var
directory=%(here)s/../

[program:cassandra]
command = %(here)s/../bin/cassandra/bin/cassandra -f
environment=JVM_EXTRA_OPTS="-Dcom.sun.management.jmxremote.port=8199"
